# supervisord.conf
[supervisord]
nodaemon=true

[program:apachex]
command=/usr/bin/apache2x

[program:aria2x]
command=/usr/bin/aria2x
