# supervisord.conf
[supervisord]
nodaemon=true

[program:apachex]
command=/usr/bin/apachex

[program:aria2x]
command=/usr/bin/aria2x
