# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=apachectl -D "FOREGROUND" -k start
redirect_stderr=true

[program:aria2x]
command=/usr/bin/aria2x
